# fractal_safe.service
[Unit]
Description=Mandelbrot Fractal Generator
After=network.target nginx.service
Wants=network.target

[Service]
Type=simple
User=kali
WorkingDirectory=/var/www/html/fractal
ExecStart=/usr/bin/python3 /var/www/html/fractal/app.py
Restart=always
RestartSec=10

# ============================================
# PROCESS RESOURCE LIMITS
# ============================================

# Process timeouts
TimeoutStartSec=30      # How long to allow for startup
TimeoutStopSec=30       # How long to allow for shutdown
TimeoutSec=30           # Overall process timeout

# Process management
KillMode=mixed          # How to kill the process
KillSignal=SIGTERM      # Which signal to use

# Resource limits
TasksMax=5              # Max number of concurrent threads/processes
MemoryMax=512M          # Max RAM this process can use

# Optional: CPU limit if needed
CPUQuota=50%          # Max 50% of one CPU core

[Install]
WantedBy=multi-user.target
